# openclaw-api/com.openclaw.api-gateway.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.openclaw.api-gateway</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/python3</string>
        <string>/Users/xiaoall/hermes-istoreos-plugin/openclaw-api/server.py</string>
        <string>--port</string>
        <string>9120</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/xiaoall/hermes-istoreos-plugin/openclaw-api</string>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <true/>
    
    <key>StandardOutPath</key>
    <string>/tmp/com.openclaw.api-gateway.log</string>
    
    <key>StandardErrorPath</key>
    <string>/tmp/com.openclaw.api-gateway.err</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>
</dict>
</plist>
